# spinnaker-monitoring-third-party/third_party/prometheus/prometheus.conf
start on filesystem or runlevel [2345]

exec /opt/PROMETHEUS_VERSION/prometheus \
  -config.file /opt/PROMETHEUS_VERSION/spinnaker-prometheus.yml \
  -storage.local.path /opt/PROMETHEUS_VERSION/data \
  2>&1 > /var/log/prometheus.log
