# spinnaker-monitoring-third-party/third_party/prometheus/prometheus.conf
start on filesystem or runlevel [2345]

exec /opt/prometheus/prometheus \
  --config.file /opt/prometheus/spinnaker-prometheus.yml \
  --storage.tsdb.path /opt/prometheus-data \
  2>&1 > /var/log/prometheus.log
