# spinnaker-monitoring-third-party/third_party/prometheus/pushgateway.conf
start on filesystem or runlevel [2345]

exec /usr/bin/pushgateway 2>&1 > /var/log/pushgateway.log
