# etc/init/spinnaker-monitoring.conf
start on filesystem or runlevel [2345]

exec /opt/spinnaker-monitoring/bin/spinnaker-monitoring.sh monitor > /var/log/spinnaker/monitoring/monitoring.log 2>&1 